# service/sshd/run
#!/bin/sh

exec </dev/null
exec 2>&1
test -d /run/sshd || mkdir /run/sshd || true
exec /usr/sbin/sshd -De
